# etc/xdg/autostart/imchooser.desktop
[Desktop Entry]
Name=IM Chooser
Comment=Choose user defined Input Method
Exec=imchooser-start
Type=Application
NoDisplay=true
